# data/rocuses/upstart/rocusagent.conf
# idmapd - NFSv4 id <-> name mapper

# rpc.idmapd is the NFSv4 ID <-> name mapping daemon.  It provides
# functionality to the NFSv4 kernel client and server, to which it
# communicates via upcalls, by translating user and group IDs to names, and
# vice versa.

description	"Rocuses Monitoring Agent"
author		"Toshiufmi Sakaguchi"

start on local-filesystems
stop on starting shutdown

respawn

exec /usr/bin/rocusagent


